# debian/pop-snapshot.prerm
#!/bin/sh
set -e

systemctl disable pop-snapshot
systemctl stop pop-snapshot

exit 0
