# debian/pop-snapshot.postinst
#!/bin/sh
set -e

systemctl enable --now pop-snapshot

exit 0
