# debian/pop-snapshot.postinst
#!/bin/sh
set -e

systemctl enable --now com.system76.PopSnapshot

exit 0
